# socketTCPpython.conf
start on started sshd
stop on runlevel [!2345]

exec /usr/bin/python ~/pyserverSocket/socketcp_gsm_py.py
respawn
